# CrewCTF2025/proxies/supervisord_proxy.conf
[supervisord]
nodaemon=true

[program:proxy]
command=/app/proxy
autostart=true
autorestart=true

[program:flask]
command=python3 /app/proxy_service.py
autostart=true
autorestart=true